# backend/deploy/systemd/mycleancenter.service
[Unit]
Description=MyCleanCenter CRM Backend
Documentation=https://github.com/mycleancenter/crm
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mycleancenter
Group=mycleancenter
WorkingDirectory=/opt/mycleancenter/current/backend
Environment=NODE_ENV=production
Environment=TZ=Europe/Berlin
Environment=DATA_DIR=/var/lib/mycleancenter
Environment=PORT=8787
Environment=HOST=0.0.0.0
Environment=FRONTEND_DIR=/opt/mycleancenter/current/dist
# CORS: Frontend und Backend laufen auf demselben Port → Same-Origin.
# Damit Zugriff per LAN-IP, Hostname und mDNS gleichermaßen funktioniert,
# erlauben wir hier alle Origins. Echte Sicherheit kommt aus Cookie-Auth +
# Lockout, nicht aus CORS.
Environment=CORS_ORIGINS=*
# Cookies dürfen NICHT als Secure laufen, solange im LAN per HTTP zugegriffen
# wird — sonst sendet der Browser sie nicht zurück und Logins wirken kaputt.
Environment=COOKIE_SECURE=false
# Google-OAuth-Redirect: bei Bedarf in der Google Cloud Console exakt eintragen.
Environment=GOOGLE_OAUTH_REDIRECT=http://mycleancenter-pi.local:8787/einstellungen/google-drive/callback
ExecStart=/usr/bin/node dist/server.js
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mycleancenter

# Sicherheits-Hardening
ProtectSystem=strict
ReadWritePaths=/var/lib/mycleancenter /opt/mycleancenter
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
